# com.foodscan3d.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
          "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  FoodScan 3D — macOS Auto-Start Service
  
  To install (auto-start on login):
    cp com.foodscan3d.server.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.foodscan3d.server.plist
  
  To uninstall:
    launchctl unload ~/Library/LaunchAgents/com.foodscan3d.server.plist
    rm ~/Library/LaunchAgents/com.foodscan3d.server.plist
  
  To start/stop manually:
    launchctl start com.foodscan3d.server
    launchctl stop  com.foodscan3d.server
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.foodscan3d.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alex/Library/Mobile Documents/com~apple~CloudDocs/A Vibe AntiGravity /3D Food from Photos/backend/venv/bin/python3</string>
    <string>/Users/alex/Library/Mobile Documents/com~apple~CloudDocs/A Vibe AntiGravity /3D Food from Photos/backend/server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alex/Library/Mobile Documents/com~apple~CloudDocs/A Vibe AntiGravity /3D Food from Photos/backend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>8000</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <!-- Start at login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart if it crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Log output -->
  <key>StandardOutPath</key>
  <string>/Users/alex/Library/Mobile Documents/com~apple~CloudDocs/A Vibe AntiGravity /3D Food from Photos/backend/server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/alex/Library/Mobile Documents/com~apple~CloudDocs/A Vibe AntiGravity /3D Food from Photos/backend/server.log</string>
</dict>
</plist>
